# Proyectos_de_Ciberseguridad/Automatizador_Servidor_ALerta_Vulnerabilidades/vuln_alerts.timer
[Unit]
Description=Ejecutar NVD Scanner cada 30 minutos

[Timer]
# Esperar 5 minutos tras el arranque para dejar que la red se estabilice
OnBootSec=5min

# Repetir cada 30 minutos después de la última ejecución
OnUnitActiveSec=30min

# Si el sistema estaba apagado, ejecutar inmediatamente al encender
Persistent=true

[Install]
WantedBy=timers.target